# ophanion-main/ophanion.service
[Unit]
Description=OPHANION - Resonant Monolith DDoS Protection for Tor Hidden Services
Documentation=https://github.com/ophanion/ophanion
After=network-online.target tor.service
Wants=network-online.target
Requires=tor.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/ophanion
ExecStart=/opt/ophanion/target/release/ophanion --config /etc/ophanion/config.toml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/ophanion /var/run/tor
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment="RUST_LOG=info"
Environment="RUST_BACKTRACE=1"

[Install]
WantedBy=multi-user.target
